# ppython/supervisord/ppython.ini
[program:ppython]
command=python3 php_python.py & ; 命令
directory=/Users/liuzhen20/work/2019/php-python/ppython/python ; 命令执行目录
stdout_logfile=/Users/liuzhen20/work/2019/php-python/ppython/log/ppython.log ;
stderr_logfile=/Users/liuzhen20/work/2019/php-python/ppython/log/ppython_error.log ;
autostart=true                ; start at supervisord start (default: true)
startsecs=10                  ; # of secs prog must stay up to be running (def. 1)
autorestart=true               ; autorestart if exited after running (def: unexpected)
stopwaitsecs=60
priority=990   ; 优先级
